# plays/instant-replay/rootfs/etc/s6-overlay/s6-rc.d/svc-autopin/run
#!/command/with-contenv bash
export IPFS_PATH=/data/ipfs

INTERVAL="${AUTOPIN_INTERVAL:-30}"
PATTERNS="${AUTOPIN_PATTERNS:-*}" 
PINNED_LOG="/data/ipfs/.autopin-state"
STAGING="/tmp/autopin"
MINIO_P="${MINIO_P:-9000}"
AUTOPIN_BUCKETS="${AUTOPIN_BUCKETS:-snapshots,static}"
MINIO_ROOT_USER="${MINIO_ROOT_USER:-minioadmin}"
MINIO_ROOT_PASSWORD="${MINIO_ROOT_PASSWORD:-minioadmin}"
ALIAS="local"
META_OBJ="metadata.json"

mkdir -p "$STAGING"
touch "$PINNED_LOG"

# Set up MinIO alias
mc alias set "$ALIAS" "http://localhost:${MINIO_P}" "$MINIO_ROOT_USER" "$MINIO_ROOT_PASSWORD" --api S3v4 >/dev/null 2>&1

echo "[autopin] Watching buckets: ${AUTOPIN_BUCKETS} every ${INTERVAL}s (patterns: ${PATTERNS})"
# Convert to array
IFS=',' read -ra BUCKETS <<< "$AUTOPIN_BUCKETS"

update_metadata() {
    local bucket="$1" objname="$2" cid="$3" etag="$4"
    local now=$(date -u +%Y-%m-%dT%H:%M:%SZ)
    local meta_local="${STAGING}/${bucket}-${META_OBJ}"

    mc cp "${ALIAS}/${bucket}/${META_OBJ}" "$meta_local" >/dev/null 2>&1 || \
        printf '{"updated_at":"","objects":{}}' > "$meta_local"

    jq --arg name "$objname" \
       --arg cid "$cid" \
       --arg etag "$etag" \
       --arg now "$now" \
       '.updated_at = $now | .objects[$name] = {ipfs_cid: $cid, etag: $etag, pinned_at: $now}' \
       "$meta_local" > "${meta_local}.tmp" && mv "${meta_local}.tmp" "$meta_local"
    mc cp "$meta_local" "${ALIAS}/${bucket}/${META_OBJ}" >/dev/null 2>&1 || \
        echo "[autopin] Warning: Failed to update metadata in ${bucket}"

    local meta_cid
    meta_cid=$(s6-setuidgid app ipfs add --nocopy --pin=true -Q "$meta_local" 2>/dev/null)
    [ -n "$meta_cid" ] && echo "[autopin] Metadata pinned -> ${meta_cid} (bucket: ${bucket})"
}

while true; do
    for bucket in "${BUCKETS[@]}"; do
        bucket=$(echo "$bucket" | xargs)
        [ -z "$bucket" ] && continue

        echo "[autopin] Scanning bucket: ${bucket}"

        # Better ls + stat approach
        mc ls --recursive "${ALIAS}/${bucket}" 2>/dev/null | awk '{print $NF}' | while IFS= read -r objname; do
            [ -z "$objname" ] && continue
            [ "$(basename "$objname")" = "$META_OBJ" ] && continue

            # Get ETag more reliably
            etag=$(mc stat "${ALIAS}/${bucket}/${objname}" 2>/dev/null | grep -i 'ETag:' | awk '{print $2}' | tr -d '"' || echo "")
            [ -z "$etag" ] && continue

            key="${bucket}/${objname}::${etag}"
            grep -qxF "$key" "$PINNED_LOG" 2>/dev/null && continue

            full_fs_path="/data/minio/${bucket}/${objname}"

            if [ ! -f "$full_fs_path" ]; then
                echo "[autopin] Warning: File not found on disk: ${full_fs_path}"
                continue
            fi

            echo "[autopin] Pinning: ${bucket}/${objname}"

            CID=$(s6-setuidgid app ipfs add --nocopy --pin=true -Q "$full_fs_path" 2>&1)

            if [[ "$CID" =~ ^(Qm[1-9A-HJ-NP-Za-km-z]{44}|bafy|bafk) ]]; then
                echo "$key" >> "$PINNED_LOG"
                echo "[autopin] Successfully pinned ${bucket}/${objname} → ${CID}"
                echo "[autopin] Gateway: http://localhost:${IPFS_GATEWAY_P:-8081}/ipfs/${CID}"

                update_metadata "$bucket" "$objname" "$CID" "$etag"
            else
                echo "[autopin] Failed to pin ${bucket}/${objname} (output: ${CID})"
            fi
        done
    done

    sleep "$INTERVAL"
done