# plays/instant-replay/rootfs/etc/s6-overlay/s6-rc.d/svc-ipfs/run
#!/command/with-contenv bash
export IPFS_PATH=/data/ipfs
export HOME=/home/app
exec s6-setuidgid app ipfs daemon --enable-gc
